# launchd/com.bta.strategy-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
<key>Label</key><string>com.bta.strategy-daily</string>
<key>ProgramArguments</key><array><string>/bin/zsh</string><string>-lc</string><string>/Volumes/Almacen/Desarrollo/bots-trading-autonomos/scripts/run_strategy_daily.sh</string></array>
<key>RunAtLoad</key><true/>
<key>StartCalendarInterval</key><dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>15</integer></dict>
<key>StandardOutPath</key><string>/Volumes/Almacen/Desarrollo/bots-trading-autonomos/runtime/logs/launchd_strategy.out.log</string>
<key>StandardErrorPath</key><string>/Volumes/Almacen/Desarrollo/bots-trading-autonomos/runtime/logs/launchd_strategy.err.log</string>
</dict></plist>
